# deploy/pyfoswiki-web.service
[Unit]
Description=PyFoswiki Web UI (FastAPI/uvicorn)
After=network.target pyfoswiki-api.service
Wants=pyfoswiki-api.service

[Service]
Type=simple
User=pyfoswiki
Group=pyfoswiki
WorkingDirectory=/opt/pyfoswiki
EnvironmentFile=/opt/pyfoswiki/.env
ExecStart=/opt/pyfoswiki/.venv/bin/uvicorn \
    webui.app:app \
    --host 127.0.0.1 \
    --port 8221 \
    --workers 2 \
    --log-level info \
    --access-log
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pyfoswiki-web

[Install]
WantedBy=multi-user.target
