# deploy/pyroswiki-api.service
[Unit]
Description=Pyroswiki API (FastAPI/uvicorn)
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=pyroswiki
Group=pyroswiki
WorkingDirectory=/opt/pyroswiki
EnvironmentFile=/opt/pyroswiki/.env
ExecStart=/opt/pyroswiki/.venv/bin/uvicorn \
    app.main:app \
    --host 127.0.0.1 \
    --port 8621 \
    --workers 2 \
    --log-level info \
    --access-log
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pyroswiki-api

[Install]
WantedBy=multi-user.target
